# server/supervisor/local.ini
[program:harvest-openmrs-local]
directory = /home/devuser/webapps/harvest-openmrs-env/harvest-openmrs
command = /home/devuser/webapps/harvest-openmrs-env/bin/uwsgi --ini /home/devuser/webapps/harvest-openmrs-env/harvest-openmrs/server/uwsgi/local.ini
autostart = true
autorestart = true
