# server/supervisor/production.ini
[program:harvest-openmrs]
command = /home/devuser/webapps/harvest-openmrs-env/bin/uwsgi --ini /home/devuser/webapps/harvest-openmrs-env/harvest-openmrs/server/uwsgi/production.ini
autostart = true
autorestart = true
